# debian/amusewiki.cron.daily
#!/bin/sh

# replaced by systemd unit files

exit 0
